# systemd/security-report.timer
[Unit]
Description=Weekly Security Report Generation

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
